# templates/local.ultra.dictation.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>             <string>local.ultra.dictation</string>
  <key>ProgramArguments</key>
  <array>
    <string>__HOME__/.local/ultra_dictation/venv/bin/python</string>
    <string>__HOME__/.local/ultra_dictation/engine.py</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>            <string>__HOME__</string>
  </dict>
  <key>StandardOutPath</key>   <string>__HOME__/.local/ultra_dictation/dictation.log</string>
  <key>StandardErrorPath</key> <string>__HOME__/.local/ultra_dictation/dictation.log</string>
  <key>RunAtLoad</key>         __RUN_AT_LOAD__
  <key>KeepAlive</key>
  <dict><key>SuccessfulExit</key><false/></dict>
  <key>ThrottleInterval</key>  <integer>5</integer>
</dict>
</plist>
